# src/systemd/rpm-ostree-toolbox-kinit.timer
[Unit]
Description=rpm-ostree-kinit timer

[Timer]
OnUnitInactiveSec=2h
Unit=rpm-ostree-kinit.service
